# recurring_costs/crontab
# START CRON JOB
0 1 * * * python3 /app/main.py
# end